# server-setup/systemd/ankaa-api.service
[Unit]
Description=Ankaa API Server (NestJS)
Documentation=https://github.com/ankaa-design/ankaa
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=kennedy
Group=ankaa
WorkingDirectory=/home/kennedy/repositories/api

# Environment
Environment=NODE_ENV=production
Environment=PORT=3030
# Node.js memory optimization - allow up to 4GB heap for smooth performance
Environment=NODE_OPTIONS="--max-old-space-size=4096"
EnvironmentFile=/home/kennedy/repositories/api/.env.production

# Node.js executable
ExecStart=/usr/bin/node dist/main.js

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Resource limits - Server has 29GB RAM, max 20 users
# Generous limits for fluid performance
LimitNOFILE=65535
MemoryMax=12G
MemorySwapMax=4G
CPUQuota=400%

# Note: Security hardening (ProtectHome, PrivateTmp, etc.) has been removed
# because it interferes with Chrome/Puppeteer used by WhatsApp integration.
# Chrome requires write access to ~/.local, ~/.config, and ~/.cache directories.

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ankaa-api

[Install]
WantedBy=multi-user.target
